# systemd/wlan0-ip.service
[Unit]
Description=WiFi Access Point IP Assignment
After=network-pre.target
Before=hostapd.service dnsmasq.service
DefaultDependencies=no

[Service]
Type=oneshot
# Read IP from hostapd config if it exists, otherwise use default
ExecStartPre=/bin/sh -c 'if [ -f /etc/hostapd/wlan0-ip.conf ]; then source /etc/hostapd/wlan0-ip.conf; else WLAN0_IP=10.0.2.1; fi'
ExecStart=/bin/sh -c 'WLAN0_IP=$(cat /etc/hostapd/wlan0-ip.conf 2>/dev/null || echo 10.0.2.1); /usr/sbin/ip addr flush dev wlan0 2>/dev/null || true; /usr/sbin/ip addr add $WLAN0_IP/24 dev wlan0 2>/dev/null || true; /usr/sbin/ip link set wlan0 up 2>/dev/null || true'
RemainAfterExit=yes

[Install]
WantedBy=sysinit.target
