# systemd/wled-manager.service
[Unit]
Description=WLED Manager - Network Management for USB Ethernet and WLED Controllers
After=network.target NetworkManager.service dnsmasq.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/wled-manager
Environment="PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin:/opt/wled-manager/venv/bin"
ExecStart=/opt/wled-manager/venv/bin/python3 /opt/wled-manager/src/app.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
